# linux_pkgs/localrepos/zoom/zoom-repo.timer
[Unit]
Description=Triggers weekly Zoom update checks
After=network.target

[Timer]
OnCalendar=15:00
OnStartupSec=20s
AccuracySec=5s
RandomizedDelaySec=5s
Persistent=true

[Install]
WantedBy=timers.target